# infra/vllm/systemd/gemma-forge-sentry.service
# GemmaForge Sentry — Gemma 4 E2B-IT, GPU 3
#
# Lightweight watchdog. Tails dmesg/journalctl on the target VM and
# classifies collateral damage between Ralph loop iterations.
#
# Port: 8070 (HTTP, OpenAI-compatible)

[Unit]
Description=GemmaForge Sentry — Gemma 4 E2B-IT (GPU 3)
Documentation=file:///data/code/gemma-forge/infra/vllm/README.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
Restart=on-failure
RestartSec=10s
TimeoutStartSec=600
TimeoutStopSec=60s
StandardOutput=append:/data/triton/logs/gemma-forge-sentry.log
StandardError=append:/data/triton/logs/gemma-forge-sentry.log

ExecStart=/data/code/gemma-forge/infra/vllm/scripts/serve.sh \
    gemma-forge-sentry \
    /weights/gemma-4-E2B-it \
    3 \
    8070 \
    --max-model-len 4096 \
    --dtype bfloat16

ExecStop=/usr/bin/docker stop --time=60 gemma-forge-sentry

[Install]
WantedBy=multi-user.target
