# infra/vllm/systemd/gemma-forge-auditor.service
# GemmaForge Auditor — Nemotron-3-Nano-30B-A3B NVFP4, PP=2, GPUs 2+3
#
# Cross-model evaluation: NVIDIA Nemotron evaluates Google Gemma's work.
# Pipeline parallelism splits layers across GPUs (not matrices), giving
# 1.2M token KV cache (120x more than Gemma's TP=2).
#
# Port: 8060 (HTTP, OpenAI-compatible)
# Tool calling: qwen3_coder parser + nano_v3 reasoning parser

[Unit]
Description=GemmaForge Auditor — Nemotron-3-Nano-30B NVFP4 (PP=2, GPUs 2+3)
Documentation=file:///data/code/gemma-forge/infra/vllm/README.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
Restart=on-failure
RestartSec=15s
TimeoutStartSec=600
TimeoutStopSec=60s
StandardOutput=append:/data/triton/logs/gemma-forge-auditor.log
StandardError=append:/data/triton/logs/gemma-forge-auditor.log

ExecStart=/data/code/gemma-forge/infra/vllm/scripts/serve-auditor.sh
ExecStop=/usr/bin/docker stop --time=60 gemma-forge-auditor

[Install]
WantedBy=multi-user.target
