# infra/vllm/systemd/gemma-forge-architect.service
# GemmaForge Architect + Worker — Gemma 4 31B-IT NVFP4, tp=2, GPUs 0+1
#
# Serves the NVFP4-quantized 31B model across GPUs 0+1 with tensor
# parallelism. Architect and Worker share this engine because they
# execute sequentially within a Ralph loop iteration.
#
# Port: 8050 (HTTP, OpenAI-compatible)
# See ADR-0015 for the model lineup rationale.

[Unit]
Description=GemmaForge Architect — Gemma 4 31B-IT NVFP4 (tp=2, GPUs 0+1)
Documentation=file:///data/code/gemma-forge/infra/vllm/README.md
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
Restart=on-failure
RestartSec=15s
TimeoutStartSec=600
TimeoutStopSec=120s
StandardOutput=append:/data/triton/logs/gemma-forge-architect.log
StandardError=append:/data/triton/logs/gemma-forge-architect.log

ExecStart=/data/code/gemma-forge/infra/vllm/scripts/serve.sh \
    gemma-forge-architect \
    /weights/Gemma-4-31B-IT-NVFP4 \
    0,1 \
    8050 \
    --tensor-parallel-size 2 \
    --max-model-len 8192

ExecStop=/usr/bin/docker stop --time=120 gemma-forge-architect

[Install]
WantedBy=multi-user.target
